# deploy/someaday-api.service.example
[Unit]
Description=SomeADay API
After=network.target

[Service]
WorkingDirectory=/srv/someaday
EnvironmentFile=/etc/someaday/someaday.env
ExecStart=/srv/someaday/.venv/bin/uvicorn services.api.main:app --host 127.0.0.1 --port 8000
Restart=on-failure

[Install]
WantedBy=multi-user.target
